# root-alarm/etc/systemd/system/alarm-pause.service
[Unit]
Description=make a pause alarm

[Service]
ExecStart=/usr/local/bin/beep pause
Type=oneshot
User=schurig
