# root-alarm/etc/systemd/system/alarm-pauseende.timer
[Unit]
Description=Timer für Mittagspause

[Timer]
AccuracySec=1
OnCalendar=Mon-Fri *-*-* 11:59:00
Unit=alarm-start.service
