# root-alarm/etc/systemd/system/alarm-start.service
[Unit]
Description=make a start alarm

[Service]
ExecStart=/usr/local/bin/beep start
Type=oneshot
User=schurig
