# root-alarm.stow/etc/systemd/system/alarm-pauseende.timer
[Unit]
Description=Timer für's Ende der Mittagspause

[Timer]
AccuracySec=1
OnCalendar=Mon-Fri *-*-* 12:29:00
Unit=alarm-start.service
