# root-email/etc/profile.d/email.sh
if [ "$HOSNAME" = "desktop" -o "$HOSTNAME" = "laptop" ]; then
	export EMAIL="holger.schurig@advantech-dlog.com"
else
	export EMAIL="holgerschurig@gmail.com"
fi
